# infra/krake_app/docker-compose/krake_api_ctrl/supervisord.conf
[supervisord]
nodeamon=true

[program:uwsgi]
command=/usr/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
